# conf/worker.conf
[supervisord]
nodaemon=true

[program:spark-slave]
command=/opt/spark/sbin/start-slave.sh spark://WORKER_IP:7077 --cores 3 --memory 9g
